# packaging/aur/slskd.service
[Unit]
Description=slskd - Soulseek client/server
Documentation=https://github.com/snapetech/slskdn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=slskd
Group=slskd
ExecStart=/usr/lib/slskd/slskd --config /etc/slskd/slskd.yml
WorkingDirectory=/var/lib/slskd
Restart=on-failure
RestartSec=10

# Security hardening (relaxed to allow symlinks to external storage)
NoNewPrivileges=yes
ProtectSystem=full
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
